# kvm-xfstests/test-appliance/files/etc/systemd/system/gce-finalize.timer
[Unit]
Description=GCE finalization timer

[Timer]
OnBootSec=12h
Unit=gce-finalize.service

[Install]
WantedBy=timers.target
