# kvm-xfstests/test-appliance/files/etc/systemd/system/gce-finalize.service
[Unit]
Description=GCE finalization

[Service]
Type=simple
ExecStart=/usr/local/lib/gce-finalize-timeout
TimeoutSec=0
StandardOutput=tty
