# test-appliance/files/etc/systemd/system/gce-repo-cleanup.service
[Unit]
Description=KCS repository cleanup

[Service]
Type=simple
ExecStart=/usr/local/lib/gce-repo-cleanup
TimeoutSec=0
StandardOutput=tty
