# test-appliance/files/etc/systemd/system/stress.service
[Unit]
Description=Stress Workload

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/lib/start-stress
